# server/deploy/zaucker-calendar.service
[Unit]
Description=Zaucker availability calendar service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/calendar
# Source nvm so the user-space (nvm) Node is on PATH (uses the `default` alias),
# then run the service. Survives Node upgrades as long as `default` is set.
ExecStart=/bin/bash -c 'export NVM_DIR=%h/.nvm; . "$NVM_DIR/nvm.sh"; exec node dist/index.js'
Environment=PORT=4317
# Bind on all interfaces so nginx on proxy-volki-01 can reach it across the
# network. The API is read-only public availability data (no secrets in
# responses); restrict with a firewall to proxy-volki-01 if you want.
Environment=HOST=0.0.0.0
Environment=CONFIG_PATH=%h/calendar/config.local.json
Environment=CACHE_PATH=%h/calendar/.cache/availability.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
