# server/deploy/zaucker-calendar.service
[Unit]
Description=Zaucker availability calendar service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/calendar
# Launch through nvm-exec so the user-space (nvm) Node is used. `nvm alias
# default` (set by `nvm install --lts`) picks the version; survives upgrades.
ExecStart=%h/.nvm/nvm-exec node dist/index.js
Environment=PORT=4317
# Bind on all interfaces so nginx on proxy-volki-01 can reach it across the
# network. The API is read-only public availability data (no secrets in
# responses); restrict with a firewall to proxy-volki-01 if you want.
Environment=HOST=0.0.0.0
Environment=CONFIG_PATH=%h/calendar/config.local.json
Environment=CACHE_PATH=%h/calendar/.cache/availability.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
