# server/deploy/zaucker-calendar.service
[Unit]
Description=Zaucker availability calendar service
After=network.target

[Service]
Type=simple
User=zaucker
WorkingDirectory=/home/zaucker/calendar
ExecStart=/usr/bin/node dist/index.js
Environment=PORT=4317
Environment=CONFIG_PATH=/home/zaucker/calendar/config.local.json
Environment=CACHE_PATH=/home/zaucker/calendar/.cache/availability.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
